# chaos_files/system_files/usr/lib/systemd/user/chaossynergy-herdr.service
[Unit]
Description=Chaossynergy — herdr terminal multiplexer
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/ptyxis -- /usr/bin/chaossynergy
Restart=no

[Install]
WantedBy=graphical-session.target